# templates/etc/init.d/Debian/carbon-cache.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-cache
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-cache
# Description: Enables Graphites carbon-cache data cache engine
### END INIT INFO

PYTHON_CMD='python -W ignore'
CARBON_DAEMON="cache"
GRAPHITE_DIR="/opt/graphite"
STOP_COUNTER=12 # 12 times 5s -> 60 secs

OPERATION="$1"
if [ $# -gt 1 ]; then
    shift
    INSTANCES=$*
else
    INSTANCES=$(grep "^\[${CARBON_DAEMON}" ${GRAPHITE_DIR}/conf/carbon.conf | cut -d \[ -f 2 | cut -d \] -f 1 | cut -d : -f 2)
fi

case "${OPERATION}" in
    start)
        for INSTANCE in ${INSTANCES}; do
            if [ "${INSTANCE}" = "${CARBON_DAEMON}" ]; then
                INSTANCE="a";
            fi;
            ${PYTHON_CMD} ${GRAPHITE_DIR}/bin/carbon-${CARBON_DAEMON}.py --instance=${INSTANCE} start
        done
        ;;
    stop)
        for INSTANCE in ${INSTANCES}; do
            if [ "${INSTANCE}" = "${CARBON_DAEMON}" ]; then
                INSTANCE="a";
            fi;
            ${PYTHON_CMD} ${GRAPHITE_DIR}/bin/carbon-${CARBON_DAEMON}.py --instance=${INSTANCE} stop
            sleep 3
            CNT=${STOP_COUNTER}
            while [ $(/usr/bin/pgrep -cf "carbon-${CARBON_DAEMON}.py --instance=${INSTANCE}") -eq 1 ]; do
                echo "waiting another 5 secs for carbon-${CARBON_DAEMON}-${INSTANCE} to stop"
                sleep 5
                if [ ${CNT} -lt 1 ]; then
                    pid=$(/usr/bin/pgrep -f "carbon-${CARBON_DAEMON}.py --instance=${INSTANCE}")
                    test -n ${pid} && kill -9 ${pid}
                    break
                fi
                CNT=$(expr ${CNT} - 1)
            done
        done
        ;;
    status)
        for INSTANCE in ${INSTANCES}; do
            if [ "${INSTANCE}" = "${CARBON_DAEMON}" ]; then
                INSTANCE="a";
            fi;
            ${PYTHON_CMD} ${GRAPHITE_DIR}/bin/carbon-${CARBON_DAEMON}.py --instance=${INSTANCE} status
        done
        ;;
    restart)
        $0 stop ${INSTANCES} && sleep 2 && $0 start ${INSTANCES}
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status} [instance instance ...]"
        exit 2
        ;;
esac

exit
