# .puppet/modules/docker/templates/etc/systemd/system/docker.socket.d/socket-overrides.conf.erb
[Socket]
SocketGroup=<%= @socket_group %>
