# .puppet/modules/prometheus/templates/prometheus.upstart.erb
# Prometheus Monitoring Framework (Upstart unit)
description "Prometheus Monitoring Framework"
start on runlevel [2345]
stop on runlevel [06]

env PROMETHEUS=<%= scope.lookupvar('prometheus::server::bin_dir') %>/prometheus
env USER=<%= scope.lookupvar('prometheus::server::user') %>
env GROUP=<%= scope.lookupvar('prometheus::server::group') %>
env DEFAULTS=/etc/default/prometheus
env RUNDIR=/var/run/prometheus
env PID_FILE=/var/run/prometheus/prometheus.pid
pre-start script
  [ -e $DEFAULTS ] && . $DEFAULTS

  mkdir -p $RUNDIR           || true
  chmod 0750 $RUNDIR         || true
  chown $USER:$GROUP $RUNDIR || true
end script

script
    # read settings like GOMAXPROCS from "/etc/default/prometheus", if available.
    [ -e $DEFAULTS ] && . $DEFAULTS

    export GOMAXPROCS=${GOMAXPROCS:-2}
    exec start-stop-daemon --chuid $USER --group $GROUP --pidfile $PID_FILE --exec $PROMETHEUS --start -- <%= @daemon_flags.join(" \\\n      ") %>
end script

respawn
respawn limit 10 10
kill timeout 10
