# modules/graphite/templates/etc/systemd/gunicorn.socket.erb
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=<%= scope.lookupvar('graphite::gunicorn_bind')[/unix:(.*)/, 1] %>

[Install]
WantedBy=sockets.target
